# clients/go/init/polda.timer
[Unit]
Description=Policy daemon service schedule

[Timer]
OnUnitActiveSec=300

[Install]
WantedBy=timers.target
